# files/system/etc/xdg/autostart/Moonlight.desktop
[Desktop Entry]
Exec=flatpak run com.moonlight_stream.Moonlight
Name=Moonlight
Type=Application
Version=1.0
